# usr/local/etc/periodic/hourly/210.clamav-unofficial-sigs-update
#! /usr/local/bin/bash
#
# /usr/local/etc/periodic/daily/340.clamav-unofficial-sigs-update

/usr/local/bin/clamav-unofficial-sigs.sh

